# assets/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

[program:php-fpm]
command=php-fpm -F
autostart=true
autorestart=true

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true

; Laravel queue worker
[program:worker-laravel]
command=bash -c 'exec php /app/artisan queue:work --sleep=3 --tries=3 --max-time=3600'
autostart=true
autorestart=true
stdout_logfile=/var/log/worker-laravel.log
stderr_logfile=/var/log/worker-laravel.log

; Inertia SSR (only if using SSR)
[program:inertia-ssr]
command=bash -c 'exec php /app/artisan inertia:start-ssr'
autostart=true
autorestart=true
stdout_logfile=/var/log/worker-inertia-ssr.log
stderr_logfile=/var/log/worker-inertia-ssr.log
